# config/system-user/bt-autoconnect.service
[Unit]
Description=Ensure BT speaker is connected and set as default sink

[Service]
Type=oneshot
ExecStart=/home/pulse/bin/bt-autoconnect.sh

[Install]
# This service is started by bt-autoconnect.timer, so this section is minimal
# The timer handles the actual enablement
WantedBy=default.target

