# config/system-user/bt-autoconnect.timer
[Unit]
Description=Periodically run bt-autoconnect

[Timer]
OnBootSec=15
OnUnitActiveSec=15

[Install]
WantedBy=default.target

