# config/system/pulse-daily-reboot.service
[Unit]
Description=Daily reboot of Pulse kiosk
Wants=network.target

[Service]
Type=simple
ExecStart=/opt/pulse-os/bin/safe-reboot.sh "daily-reboot"
