# config/system/pulse-hw-watchdog.service
[Unit]
Description=Pulse Hardware Watchdog (kernel error monitor)
After=multi-user.target

[Service]
Type=simple
ExecStart=/opt/pulse-os/bin/pulse-hw-watchdog.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
