# config/system/pulse-kiosk-mqtt.service
[Unit]
Description=Pulse Kiosk MQTT Navigation Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/pulse-os/bin/kiosk-mqtt-wrapper.sh
Restart=always
RestartSec=3
User=pulse
WorkingDirectory=/opt/pulse-os
Environment=PYTHONPATH=/opt/pulse-os
Environment=XDG_RUNTIME_DIR=/run/user/%%U
PassEnvironment=MQTT_HOST MQTT_PORT PULSE_URL
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

