# config/system/pulse-daily-reboot.service
[Unit]
Description=Daily reboot of Pulse kiosk
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/systemctl reboot
