# config/system/pulse-backlight-sun.service
[Unit]
Description=Pulse day/night auto-adjustment (screen brightness)
After=multi-user.target

[Service]
Type=simple
User=pulse
WorkingDirectory=/opt/pulse-os
Environment=PYTHONPATH=/opt/pulse-os
ExecStart=/usr/bin/env python3 /home/pulse/bin/pulse-backlight-sun.py
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
