# config/system/syslog-ng.service.d/override.conf
[Service]
Environment=REMOTE_LOG_HOST=${PULSE_REMOTE_LOG_HOST}
Environment=REMOTE_LOG_PORT=${PULSE_REMOTE_LOG_PORT}
